# .conf/conf.d/celery_ppomppu_worker.conf
[program:ppomppuworker]
command=/usr/local/bin/celery -A conf worker -l info
directory=/srv/app/ppomppu
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=600
