# .conf/conf.d/celery_ppomppu_beat.conf
[program:ppomppubeat]
command=/usr/local/bin/celery -A conf beat -l info
directory=/srv/app/ppomppu
startsecs=10
stopwaitsecs=600
autostart=true
